# smc-backend.service
[Unit]
Description=SMC Dashboard Backend API
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/smc-backend
Environment=NODE_ENV=production
EnvironmentFile=/opt/smc-backend/.env
ExecStart=/usr/bin/node /opt/smc-backend/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smc-backend

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/smc-backend/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target

